# sv/runsvdir/finish
#!/bin/sh
sv force-shutdown /etc/service/* >/dev/null 2>&1
sv x .
